# supervisord.conf
[supervisord]
nodaemon=true

[program:dmserver]
command=/bin/bash -c "/opt/dmdbms/bin/dmserver /opt/dmdbms/DATA/DAMENG/dm.ini -noconsole"

[program:sshd]
command=/usr/sbin/sshd -D
